# fix-easytier.service
[Unit]
Description=Fix EasyTier IP Forwarding and NAT rules
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix_easytier.sh
User=root
Restart=no
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=fix-easytier

[Install]
WantedBy=multi-user.target
